# templates/default/sc@.service.erb
[Unit]
Description=Sauce Connect worker service on port %i
PartOf=sc.service
ReloadPropagatedFrom=sc.service

[Service]
Type=simple
User=nobody
Group=nobody
WorkingDirectory=/tmp
LimitNOFILE=8192
Restart=always

# Set those to match your Saucelabs credentials
Environment=SAUCE_USERNAME=<%= @sauce_username %>
Environment=SAUCE_ACCESS_KEY=<%= @sauce_access_key %>

ExecStart = /usr/local/bin/sc \
<% @sc_cmd_options.each do |opt| %>
  <%= opt %> \
<% end %>

[Install]
WantedBy=multi-user.target
